# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/root/mysql_startup.sh
